# build/COPY_ROOT_1/etc/supervisor/supervisord/conf.d/comfyui-conda.conf
[program:comfyui-conda]
command=/opt/ai-dock/bin/supervisor-comfyui-conda.sh
environment=
    PATH="/opt/miniconda/bin:%(ENV_PATH)s",
    CONDA_DEFAULT_ENV="comfyui"
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
user=root
priority=20
