# build/supervisord.ini
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root --ServerApp.token="" --ServerApp.password="" --ServerApp.allow_origin="*"
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/jupyter.log
stderr_logfile=/var/log/jupyter.log

[program:tensorboard]
command=/usr/local/bin/start-tensorboard.sh
directory=/workspace
autostart=%(ENV_TENSORBOARD_AUTOSTART)s
autorestart=true
stdout_logfile=/var/log/tensorboard.log
stderr_logfile=/var/log/tensorboard.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
user=root
startsecs=15
startretries=5
exitcodes=0
stopsignal=TERM
stopwaitsecs=10

[program:filebrowser]
command=filebrowser -a 0.0.0.0 -p 8080 -r /workspace
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.log

[program:infinite-browser]
command=python app.py --host=0.0.0.0 --port=8188
directory=/opt/infinite-browser
autostart=true
autorestart=true
stdout_logfile=/var/log/infinite-browser.log
stderr_logfile=/var/log/infinite-browser.log
