# build/supervisord.ini
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
silent=false

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:jupyter]
command=/bin/bash -c "source /opt/miniconda/bin/activate diffusion-pipe && jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root --ServerApp.token='' --ServerApp.password='' --ServerApp.allow_origin='*'"
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/jupyter.log
stderr_logfile=/var/log/jupyter.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=PATH="/opt/miniconda/envs/diffusion-pipe/bin:/opt/miniconda/bin:%(ENV_PATH)s"
priority=200
startsecs=5
startretries=3

[program:tensorboard]
command=/usr/local/bin/start-tensorboard.sh
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/tensorboard.log
stderr_logfile=/var/log/tensorboard.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=PATH="/opt/miniconda/envs/diffusion-pipe/bin:/opt/miniconda/bin:%(ENV_PATH)s",CONDA_DEFAULT_ENV="diffusion-pipe",CONDA_PREFIX="/opt/miniconda/envs/diffusion-pipe"
user=root
startsecs=15
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
priority=300

[program:filebrowser]
command=filebrowser -a 0.0.0.0 -p 8080 -r /workspace
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
redirect_stderr=false
priority=400

[program:infinite-browser]
command=/bin/bash -c "source /opt/miniconda/bin/activate diffusion-pipe && python app.py --host=0.0.0.0 --port=8188"
directory=/opt/infinite-browser
autostart=true
autorestart=true
stdout_logfile=/var/log/infinite-browser.log
stderr_logfile=/var/log/infinite-browser.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=PATH="/opt/miniconda/envs/diffusion-pipe/bin:/opt/miniconda/bin:%(ENV_PATH)s"
priority=500
